# tests/vps/sing-box-yg-health.service
[Unit]
Description=Sing-box-yg health sample

[Service]
Type=oneshot
ExecStart=/usr/local/lib/sing-box-yg/monitor-health.sh
