# tests/vps/sing-box-yg-health.service
[Unit]
Description=sing-box-yg validated health repair
After=network-online.target sing-box.service

[Service]
Type=oneshot
UMask=0077
ExecStart=/usr/local/lib/sing-box-yg/sb-doctor.sh --repair
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/etc/s-box /var/lib/sing-box-yg /run
